# supervisor.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=info

[program:live-server]
command=python -m uvicorn live.server:app --host 0.0.0.0 --port 8080 --log-level info
directory=/app
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:graph-api]
command=python -m uvicorn graph_api:app --host 0.0.0.0 --port 9000 --log-level info
directory=/app
autostart=true
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:root-agent]
command=python -m uvicorn agents.adk.root_agent.agent:a2a_app --host 0.0.0.0 --port 8000 --log-level info
directory=/app
autostart=true
autorestart=true
priority=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:caclub-agent]
command=python -m uvicorn agents.adk.caclub_a2a.agent:a2a_app --host 0.0.0.0 --port 8001 --log-level info
directory=/app
autostart=true
autorestart=true
priority=40
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:taxtmi-agent]
command=python -m uvicorn agents.adk.taxtmi_a2a.agent:a2a_app --host 0.0.0.0 --port 8002 --log-level info
directory=/app
autostart=true
autorestart=true
priority=40
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:turbotax-agent]
command=python -m uvicorn agents.adk.turbotax_a2a.agent:a2a_app --host 0.0.0.0 --port 8003 --log-level info
directory=/app
autostart=true
autorestart=true
priority=40
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:taxprofblog-agent]
command=python -m uvicorn agents.adk.taxprofblog_a2a.agent:a2a_app --host 0.0.0.0 --port 8004 --log-level info
directory=/app
autostart=true
autorestart=true
priority=40
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
